# sidechain/solfunmeme-crawl.timer
[Unit]
Description=SOLFUNMEME Daily Crawl Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
